# tftp-proxy.service
[Unit]
Description=tftp-proxy

[Service]
ExecStart=/usr/bin/tftp-proxy -url=http://zerotouch.example.com

[Install]
WantedBy=multi-user.target